# templates/etc/profile.d/kafka-consumer.sh.j2
export PATH=$PATH:/opt/kafka/bin
export KAFKA_BROKER_CONNECT="{{kafka_broker_connect}}"
